# prniolet-server.service
[Unit]
Description=Python ESCPOS Server

[Service]
Type=simple
ExecStart=/usr/local/bin/prniolet-server

[Install]
WantedBy=multi-user.target
