# extra/prniotlet-server@.service
[Unit]
Description=Python ESCPOS Server

[Service]
Type=simple
ExecStart=/usr/local/bin/prniotlet-server

[Install]
WantedBy=multi-user.target
